# IONITX-A-U/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

# allow wake on lan
ethtool -s eth0 wol g

# fix for WDEARS20 head parking one every 8 sec (what were they thinking)
hdparm -S 242 /dev/sdb
